# supervisor/storm-daemon.conf
[program:storm-%daemon%]
command=/usr/bin/storm %daemon%
directory=/home/storm
user=storm
autostart=true
autorestart=true
startsecs=20
redirect_stderr=true